# docker/radius-acct/supervisord.conf
[supervisord]
nodaemon=true

[program:certmanager]
command=/usr/local/pf/sbin/pfcertmanager

[program:radius-acct]
command=/usr/sbin/freeradius -d /usr/local/pf/raddb -n acct -fm
