# packetfence.rsyslog-drop-in.service
[Service]
ExecStartPre=/usr/local/pf/bin/pfcmd generatesyslogconfig
