# packetfence.logrotate-drop-in.service
[Service]
ReadWritePaths=/usr/local/pf/logs/
ReadWritePaths=/usr/local/pf/var/
ReadWritePaths=/usr/local/fingerbank/logs/

