# packetfence-pfconnector-remote.logrotate-drop-in.service
[Service]
ReadWritePaths=/usr/local/pfconnector-remote/logs/

