# conf/systemd/packetfence-base.slice
[Unit]
Description=The PacketFence base resource control slice
Before=slices.target

