# conf/systemd/packetfence-radiusd-eduroam.service
[Unit]
Description=PacketFence FreeRADIUS multi-protocol EDUROAM authentication server
Documentation=man:radiusd(8) man:radiusd.conf(5) http://wiki.freeradius.org/ http://networkradius.com/doc/
Wants=packetfence-httpd.aaa.service 
After=packetfence-radiusd-acct.service
Wants=packetfence-base.target packetfence-config.service packetfence-iptables.service
After=packetfence-base.target packetfence-config.service packetfence-iptables.service 

[Service]
Type=notify
WatchdogSec=10
StartLimitBurst=3
StartLimitInterval=10
PIDFile=/usr/local/pf/var/run/radiusd-eduroam.pid
ExecStartPre=/usr/local/pf/bin/pfcmd service radiusd generateconfig
ExecStartPre=/usr/sbin/radiusd -d /usr/local/pf/raddb  -n eduroam -Cxm -lstdout
ExecStart=/usr/sbin/radiusd -d /usr/local/pf/raddb  -n eduroam -fm 
Restart=on-failure
Slice=packetfence.slice

[Install]
WantedBy=packetfence.target
