# conf/systemd/packetfence.slice
[Unit]
Description=The PacketFence resource control slice
Before=slices.target

