# conf/systemd/packetfence-tracking-config.path
[Unit]
Description=PacketFence Configuration Change Tracking
Documentation=http://packetfence.org

[Path]
PathChanged=/usr/local/pf/conf/

[Install]
WantedBy=packetfence-base.target
