# docker/radius-auth/supervisord.conf
[supervisord]
nodaemon=true

[program:certmanager]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/local/pf/sbin/pfcertmanager

[program:radius-auth]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/sbin/freeradius -d /usr/local/pf/raddb -n auth -fm
