# docker/radius-auth/supervisord.conf
[supervisord]
nodaemon=true

[program:certmanager]
command=/usr/local/pf/sbin/pfcertmanager

[program:radius-auth]
command=/usr/sbin/freeradius -d /usr/local/pf/raddb -n auth -fm
