# addons/pfconnector/systemd/packetfence-pfconnector-remote.service
# Copyright (C) Inverse inc.
[Unit]
Description=PacketFence Connector Client (Remote)
Wants=packetfence-base.target packetfence-config.service packetfence-firewalld.service
After=packetfence-base.target packetfence-config.service packetfence-firewalld.service

[Service]
LimitNOFILE=infinity
EnvironmentFile=/usr/local/pfconnector-remote/conf/pfconnector-client.env
ExecStart=/usr/local/pfconnector-remote/bin/pfconnector client ENV 127.0.0.1:22226:127.0.0.1:22226
Restart=always

[Install]
WantedBy=multi-user.target
