# docker/radius-load-balancer/supervisord.conf
[supervisord]
nodaemon=true

[program:certmanager]
command=/usr/local/pf/sbin/pfcertmanager

[program:radius-load-balancer]
command=/usr/sbin/freeradius -d /usr/local/pf/raddb -n load_balancer -fm
